# ledspi-service.in
### BEGIN INIT INFO
# Provides: spileds-server
# Required-Start: $local_fs $remote_fs
# Required-Stop:
# X-Start-Before:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Starts the SPIleds server
### END INIT INFO
 
SPILEDS=LEDSPI_PATH/run-ledspi
SPILEDS_ARGS=""
PIDFILE=/var/run/spileds-server.pid
case "$1" in
  start)
    echo -n "Starting SPIleds Server: spileds-server"
    start-stop-daemon --start --quiet --pidfile $PIDFILE --make-pidfile --background --exec $SPILEDS -- $SPILEDS_ARGS
    echo "."
    ;;
  stop)
    echo -n "Stopping SPIleds Server: spileds-server"
    start-stop-daemon --stop --quiet --pidfile $PIDFILE
    echo "."
    ;;
  restart)
    $0 stop
    $0 start
    ;;
  *)
        echo "Usage: /etc/init.d/xvfb {start|stop|restart}"
        exit 1
esac
 
exit 0
